# systemd/claude-remote.service
[Unit]
Description=Persistent Claude Code remote-control session
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/cc-remote
WorkingDirectory=%h
Restart=on-failure
RestartSec=10
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
